# rwifi.service
[Service]
ExecStart=/bin/rwifi.sh


[Install]
WantedBy=multi-user.target
